# update-motd.d/00-header
#!/usr/bin/tclsh
#
#    00-header - Draw system header
#    Copyright (C) 2017 Malin von Matern
#
#    Authors: Malin von Matern
#

# ---------
# VARIABLES
# ---------

set dist(id) [exec -- lsb_release -is]
set dist(release) [exec -- lsb_release -rs]
set dist(codename) [exec -- lsb_release -cs]
set system(os) [exec -- uname -o]
set system(kernel) [exec -- uname -r]
set system(hw) [exec -- uname -m]

# ------------
# PRINT OUTPUT
# ------------

puts "                                             "
puts "           \033\[01;30m, 'b\033\[0m                              "
puts "          \033\[01;30m,\$  `Y.\033\[0m                            "
puts "          \033\[01;30m;P   \"Y.\033\[0m                           "
puts "         \033\[01;30m,\$\"    Y\$:     _\033\[0m                    "
puts "         \033\[01;30m,\$b    \$\$P. ,g\$\"´\033\[0m                   "
puts "         \033\[01;30m:\$\$.   g\$\$\";d\"\033\[0m                      "
puts "         \033\[01;30m`\$\$;,:´d\$\$\$\$P      ____\033\[0m             "
puts "          \033\[01;30mY\$\$´_d\$\$\$\$\$b:,gt\$\$\$\$\$\$\$tgg.\033\[0m        "
puts "          \033\[01;30m:\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$Y\"\"\"\"`Y\$\$p_\033\[0m     "
puts "         \033\[01;30m,\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$Y_           \"\033\[0m    "
puts "         \033\[01;30m\"Y\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$\$pm_\033\[0m            "
puts "          \033\[01;30m\$\$\$\$\"\$\$P\"\"Y\$\$\$\$\$\$\$\$P.  `\033\[0m           "
puts "         \033\[01;30md\$\$P\033\[01;36mo\033\[0m\033\[01;30m,\$b, .\$\$\$\$\$\$\$\$\$\$\$g.\033\[0m            "
puts "        \033\[01;30m,\$\$\$\$\$\$\$P\" g\$\$P\"\"\$\$b\"\"\"Y\$b,\033\[0m          "
puts "       \033\[01;30m,d\$\$\$\$\$P'   \"\$\$   `\$P     `Y,\033\[0m         "
puts "      \033\[01;30m.\$\$\$\$\$\$P  .  d\$\$,_  \"        '\033\[0m         "
puts "      \033\[01;30m:\$P\"\"Yp' .\" ,\$\$\$\"\033\[0m       \033\[1;01;36m,            \`##\033\[0m    "
puts "      \033\[01;30m`Y\$P,\"  g\" ,d\$\$\$P\"\033\[0m    \033\[1;01;36m.##.  ,,,, ,,,  ##\033\[0m    \033\[1mPOWERED BY\033\[0m"
puts "        \033\[01;30m\"\$Y  g\$ ,\$P\"\033\[0m         \033\[1;01;36m##    '#.  #   ##\033\[0m    $dist(id) $dist(release) \"$dist(codename)\""
puts "          \033\[01;30m`,g\$:\033\[0m              \033\[1;01;36m##     '#.#    ##\033\[0m    $system(os) $system(kernel)"
puts "     \033\[01;30m_____,gg'\033\[0m               \033\[1;01;36m'#,'    '#    .##,\033\[0m   $system(hw)"
puts "  \033\[01;30m,P\"\"\"\"\";\$;'\033\[0m                     \033\[1;01;36m,, #\033\[0m            "
puts "  \033\[01;30m;   ,d\"'\033\[0m                         \033\[1;01;36m''\033\[0m             "
puts "      \033\[01;30m\".\033\[0m                    ================================================"
